# onair.service
[Unit]
Description=OnAir service
After=network.target

[Service]
ExecStart=/home/ricci/onair/onair.py
